# rootfs/etc/s6-overlay/s6-rc.d/aloha/run
#!/command/with-contenv bash

export ALOHA_DATA_DIR=/data/aloha
export ALOHA_HA_CONFIG_DIR=/data/homeassistant
export HOMEASSISTANT_CONFIG=/data/homeassistant

mkdir -p /data/aloha /data/homeassistant

echo "[Aloha] Waiting for Home Assistant..."
_waited=0
until curl -sf http://localhost:8123/api/ >/dev/null 2>&1; do
    sleep 2
    _waited=$((_waited + 2))
    if [ $_waited -ge 120 ]; then
        echo "[Aloha] WARNING: Home Assistant did not respond within 120s, proceeding anyway"
        break
    fi
done

echo "[Aloha] Starting on port 7123..."
exec python3 -m aloha
